# service/configs/include/workers.ini
[program:beanstalkd]
command=beanstalkd -b /data/beanstalkd/
autostart=true
autorestart=true
stdout_logfile=/data/logs/beanstalkd.log ; (main log file;default $CWD/supervisord.log)
stdout_logfile_maxbytes=1MB

[program:worker_sms]
command=php /data/www/dachuwang/web/applications/worker/index.php worker_sms run
priortity=1
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/data/logs/worker_sms.log ; (main log file;default $CWD/supervisord.log)
stdout_logfile_maxbytes=100MB
